# Scheduler/crontab

0 12 * * * sh /app/scheduler_init.sh Parserdata-container
30 12 * * * sh /app/scheduler_init.sh Connector-container